# live/config/includes.chroot/etc/systemd/system/adduser.service
[Unit]
Description=AddUser Ceubos

[Service]
Type=oneshot
ExecStart=/usr/local/bin/adduser.sh

[Install]
WantedBy=multi-user.target


